# layout/Library/LaunchDaemons/com.procmonrootless.procmond.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.procmonrootless.procmond</string>
	<key>ProgramArguments</key>
	<array>
		<string>/var/jb/usr/libexec/procmond</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>2</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Umask</key>
	<integer>18</integer>
	<key>StandardOutPath</key>
	<string>/var/jb/var/log/procmond.launchd.out</string>
	<key>StandardErrorPath</key>
	<string>/var/jb/var/log/procmond.launchd.err</string>
</dict>
</plist>
